# deploy/repo-watch-api.cron
* * * * * www-data php /var/www/repo-watch-api/current/artisan schedule:run >> /var/log/repo-watch-api-scheduler.log 2>&1
